# system_files/etc/systemd/system/resize.service
[Unit]
Description=Resizes the file system on the first boot

[Service]
ExecStart=/bin/bash /boot/resize.sh

[Install]
WantedBy=default.target
